# cloudshell2-fan.service
[Unit]
Description=Switches the fan on/off depending on the temperatures

[Service]
Type=simple
User=root
ExecStart=/usr/bin/cloudshell2-fan
Restart=on-abort

[Install]
WantedBy=multi-user.target
